# etc/cron.d/tmpwatch-cron
*/5 * * * * root /usr/bin/tmpwatch -amv 1 /tmp
